# luci-app-happynet/root/etc/uci-defaults/luci-happynet
#!/bin/sh

uci -q batch <<-EOF >/dev/null
	delete ucitrack.@happynet[-1]
	add ucitrack happynet
	set ucitrack.@happynet[-1].init=happynet
	commit ucitrack
EOF

rm -f /tmp/luci-indexcache
exit 0
